# .config/global-files/etc/systemd/system/issue.service
[Unit]
Description=Generates /etc/issue from fortune

[Service]
Type=oneshot
ExecStart=/usr/sbin/issue

[Install]
WantedBy=multi-user.target
